# opt/profile-shutdown.service
[Unit]
Description="cURL headnode with Flight Profile remove request"
DefaultDependencies=no
Requires=network.target
Before=shutdown.target
After=sshd.service flight-slurmd.service nfs-server.service

[Service]
Type=oneshot
RemainAfterExit=true
ExecStop=/root/shutdown.sh

[Install]
WantedBy=halt.target shutdown.target
